# configurations/crontab
# crontab -e
*/30 * * * * /home/captain/smart_boat/scripts/update_kindle.sh
* * * * * /home/captain/smart_boat/scripts/sensors-influxdb.py