# scripts/ve.direct-python/lib/systemd/system/vedirect.service
[Unit]
Description=VE.Direct serial monitor
After=influxd.service

[Service]
Type=simple
User=captain
Group=captain
WorkingDirectory=/home/captain/smart_boat/scripts/ve.direct-python
EnvironmentFile=/home/captain/smart_boat/scripts/ve.direct-python/.env.installed
ExecStart=/usr/bin/python3 -u -m vedirect --influx=${INFLUX_HOST} --database=${INFLUX_DB} --port=${SERIAL_PORT}
StandardOutput=syslog
StandardError=syslog
Nice=-20
RestartSec=60s
Restart=always

[Install]
WantedBy=multi-user.target